# scripts/com.igniteiq.export-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.igniteiq.export-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/matthewlawler/scripts/igniteiq/watch-exports.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/matthewlawler/Google Drive/My Drive/IgniteIQ/Claude Design Exports</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/igniteiq-export-watcher.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/igniteiq-export-watcher.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/matthewlawler</string>
    </dict>
</dict>
</plist>
